# configs/audit.rules
# ============================================================
# Audit Rules – Server Hardening Lab
# ============================================================

-D

-b 8192

-f 1

-w /etc/passwd -p wa -k user_modification
-w /etc/shadow -p wa -k user_modification
-w /etc/group  -p wa -k user_modification

-w /etc/ssh/sshd_config -p wa -k ssh_config_change

-w /etc/sudoers -p wa -k sudoers_change
-w /etc/sudoers.d/ -p wa -k sudoers_change

-w /etc/cron.d/       -p wa -k cron_change
-w /etc/crontab       -p wa -k cron_change
-w /var/spool/cron/   -p wa -k cron_change

-a always,exit -F arch=b64 -S execve -F euid=0 -k root_commands

-w /bin/su   -p x -k privilege_escalation
-w /usr/bin/sudo -p x -k privilege_escalation

-w /etc/hosts         -p wa -k network_change
-w /etc/network/      -p wa -k network_change
-w /etc/netplan/      -p wa -k network_change

-w /var/log/faillog -p wa -k login_failures
-w /var/log/lastlog -p wa -k login_info
